# templates/fpm/fpm-kill-pool.timer.erb
[Unit]
Description=Kill fpm units without children

[Timer]
OnCalendar=<%= @on_calendar %>

[Install]
WantedBy=timers.target

